# files/etc/service/phd/run
#!/bin/bash
exec 2>&1
exec chpst -u app:app /usr/local/bin/phd-wrapper

